# sshd.conf
[program:sshd]
directory=/
command=/usr/sbin/sshd -D -E /var/log/%(program_name)s.log
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
